# install/etc/s6/services/06-rsyslog/run
#!/usr/bin/with-contenv bash

source /assets/functions/*
PROCESS_NAME="rsyslog"
check_initialized
liftoff

print_info "Starting Syslog"
silent exec rsyslogd -n
