# install/etc/s6/services/05-fail2ban/run
#!/usr/bin/with-contenv bash

source /assets/functions/*
PROCESS_NAME="fail2ban"
check_initialized
liftoff

print_info "Starting fail2ban"
silent exec /usr/bin/fail2ban-client -f start
